# pkg/etc/systemd/system/wearablepie-gpsd.service
[Unit]
Description=GPS Daemon Starter

[Service]
Type=simple
ExecStart=/usr/share/wearablepie/scripts/gpsd

[Install]
WantedBy=multi-user.target
